# api/cron/metagraph.cron
 */10 * * * * python3 /home/ec2-user/cronjobs/metagraph.py
