# api/cron/metagraph.cron
 */10 * * * * poetry run metagraph:fetch
